# promtail.service
[Unit]
Description=Promtail service for log aggregation
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/adminstd/gitpr
ExecStart=/usr/local/bin/promtail --config.file=/etc/promtail-local-config.yaml
Restart=always
RestartSec=10

# Environment variables
Environment=LOKI_URL=http://localhost:3100
Environment=PROMTAIL_PORT=9080

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=promtail

[Install]
WantedBy=multi-user.target

